# pi_software/deploy/fieldcore-kiosk.service
[Unit]
Description=FieldCore Chromium Kiosk
After=fieldcore-web.service graphical.target
Wants=graphical.target

[Service]
Type=simple
User=wex
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/wex/.Xauthority
ExecStart=/home/wex/app/pi_software/deploy/kiosk.sh
Restart=on-failure
RestartSec=10

[Install]
WantedBy=graphical.target
