# pi_software/deploy/fieldcore-lora.service
[Unit]
Description=FieldCore LoRa Packet Listener
After=fieldcore-web.service

[Service]
Type=simple
User=wex
WorkingDirectory=/home/wex/app/pi_software
Environment=FIELDCORE_DB=/home/wex/app/backend/sensors.db
ExecStart=/home/wex/env/bin/python3 -m services.lora_listener
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
